# deploy/strongpassword-mcp.service
[Unit]
Description=StrongPassword read-only MCP server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opc
WorkingDirectory=/var/www/strongpassword
Environment=HOST=127.0.0.1
Environment=MCP_PORT=3003
ExecStart=/usr/bin/node /var/www/strongpassword/mcp-server.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
